# scripts/wsmq.service
[Unit]
Description=WSMQ Service

[Service]
Type=simple
ExecStart=/usr/bin/wsmq start --port 8080

[Install]
WantedBy=multi-user.target